# ansible/roles/openbao/templates/openbao-seal-check.timer.j2
[Unit]
Description=OpenBao Seal Status Check Timer
Documentation=man:systemd.timer(5)

[Timer]
# Check every 15 minutes
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
