# ansible/roles/node-exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter - System Metrics Collector
Documentation=https://github.com/prometheus/node_exporter
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ node_exporter_user }}
Group={{ node_exporter_group }}

# Security hardening
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
PrivateDevices=false  # Need access to devices for metrics
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
ReadOnlyPaths=/

# Allow write access to textfile directory
ReadWritePaths={{ node_exporter_textfile_dir }}

# Environment
{% if node_exporter_env_vars %}
{% for key, value in node_exporter_env_vars.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
{% endif %}

# Node Exporter command
ExecStart={{ node_exporter_install_dir }}/node_exporter \
  --web.listen-address={{ node_exporter_web_listen_address }} \
{% if node_exporter_web_basic_auth_enabled or node_exporter_tls_enabled %}
  --web.config.file=/etc/node_exporter/web-config.yml \
{% endif %}
{% for collector in node_exporter_enabled_collectors %}
{% if collector is mapping %}
{% for key, value in collector.items() %}
  --collector.{{ key }} \
{% if value is mapping %}
{% for opt_key, opt_value in value.items() %}
  --collector.{{ key }}.{{ opt_key }}={{ opt_value }} \
{% endfor %}
{% endif %}
{% endfor %}
{% else %}
  --collector.{{ collector }} \
{% endif %}
{% endfor %}
{% for collector in node_exporter_disabled_collectors %}
  --no-collector.{{ collector }} \
{% endfor %}
  {{ node_exporter_extra_args }}

ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart={{ node_exporter_systemd_restart }}
RestartSec={{ node_exporter_systemd_restart_sec }}

# Resource limits
LimitNOFILE={{ node_exporter_systemd_limit_nofile }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=node_exporter

[Install]
WantedBy=multi-user.target
