# ansible/roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus - Monitoring System and Time Series Database
Documentation=https://prometheus.io/docs/
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_group }}

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

# Allow write access to data directory
ReadWritePaths={{ prometheus_data_dir }}
ReadWritePaths={{ prometheus_log_dir }}

# Environment
{% if prometheus_env_vars %}
{% for key, value in prometheus_env_vars.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
{% endif %}

# Prometheus command
ExecStart={{ prometheus_install_dir }}/prometheus \
  --config.file={{ prometheus_config_dir }}/prometheus.yml \
  --storage.tsdb.path={{ prometheus_data_dir }} \
  --storage.tsdb.retention.time={{ prometheus_storage_tsdb_retention_time }} \
{% if prometheus_storage_tsdb_retention_size != '0' %}
  --storage.tsdb.retention.size={{ prometheus_storage_tsdb_retention_size }} \
{% endif %}
  --web.listen-address={{ prometheus_web_listen_address }} \
  --web.external-url={{ prometheus_web_external_url }} \
  --web.console.templates={{ prometheus_config_dir }}/consoles \
  --web.console.libraries={{ prometheus_config_dir }}/console_libraries \
{% if prometheus_web_basic_auth_enabled or prometheus_tls_enabled %}
  --web.config.file={{ prometheus_config_dir }}/web.yml \
{% endif %}
  --log.level=info \
  {{ prometheus_extra_args }}

ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart={{ prometheus_systemd_restart }}
RestartSec={{ prometheus_systemd_restart_sec }}

# Resource limits
LimitNOFILE={{ prometheus_systemd_limit_nofile }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prometheus

[Install]
WantedBy=multi-user.target
