# ansible/roles/security_hardening/templates/run-shm.mount.j2
# Systemd Mount Unit for /run/shm
# Managed by Ansible - security_hardening role
# Location: /etc/systemd/system/run-shm.mount
#
# Purpose: Mount /run/shm (shared memory) as tmpfs with security hardening
# - noexec: Prevent execution of binaries (common attack vector)
# - nosuid: Prevent setuid/setgid bits
# - defaults: Standard mount options

[Unit]
Description=Secure Shared Memory (/run/shm)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=tmpfs
Where=/run/shm
Type=tmpfs
Options=defaults,noexec,nosuid

[Install]
WantedBy=local-fs.target
