# ansible/roles/loki/templates/loki.service.j2
[Unit]
Description=Loki log aggregation system
Documentation=https://grafana.com/docs/loki/latest/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ loki_user }}
Group={{ loki_group }}

ExecStart=/usr/bin/loki \
  -config.file={{ loki_config_dir }}/loki.yml

ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ loki_data_dir }} {{ loki_log_dir }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loki

[Install]
WantedBy=multi-user.target
